# ssh/etc/init.d/factory_stop
#!/bin/sh /etc/rc.common

START=50

start() {
	true
}